# Docker/supervisord.conf
[supervisord]
nodaemon=true
[program:ping1]
command=/bin/ping localhost
[program:ping2]
command=/bin/ping wp.pl
[program:ping3]
command=/bin/ping 8.8.8.8
[program:sshd]
command=/usr/sbin/sshd -D
